# crates/witness-generator-spec-cli/systemd/witness-exporter.timer
[Unit]
Description=Run stateless input exporter periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
