# crates/witness-generator-spec-cli/systemd/witness-collector.service
[Unit]
Description=Stateless input collector
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=stateless-inputs
Group=stateless-inputs
Environment=RUST_LOG=info
ExecStart=/usr/local/bin/witness-generator-spec-cli collect --config /etc/witness-generator-spec-cli/glamsterdam-devnet-5.toml
Restart=always
RestartSec=10

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/data/stateless-inputs

[Install]
WantedBy=multi-user.target
