# crates/witness-generator-spec-cli/systemd/witness-exporter.service
[Unit]
Description=Export stateless input batches

[Service]
Type=oneshot
User=stateless-inputs
Group=stateless-inputs
Environment=RUST_LOG=info
ExecStart=/usr/local/bin/witness-generator-spec-cli export --config /etc/witness-generator-spec-cli/glamsterdam-devnet-5.toml

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/data/stateless-inputs
